# config/docker/php/crontab
* * * * * php /app/cron.php --command=process_events_finished >> /app/cron.log 2>&1